# debian/postinst
#!/bin/sh
set -e

# Enable and start nvfan service
systemctl daemon-reload
systemctl enable nvfan.service
