# debian/prerm
#!/bin/sh
set -e

# Stop and disable nvfan service
systemctl disable --now nvfan.service 2>/dev/null || true
systemctl daemon-reload
